# deploy/systemd/foreman-quant-review.timer
[Unit]
Description=Run Foreman Quant Review Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
